# config/systemd/smoker-backend.service
[Unit]
Description=CherrySmoked Flask Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=juliusjj25
WorkingDirectory=/home/juliusjj25/SM3K/backend
ExecStart=/home/juliusjj25/SM3K/backend/start.sh
Environment="PATH=/home/juliusjj25/SM3K/backend/venv/bin"
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
